# machines/wsa/fsws-arch-config/systemd/user/nexus-orange-rgb.service
[Unit]
Description=Force static orange RGB on login
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 5
ExecStart=%h/.local/bin/set-orange-rgb.sh ff8800
RemainAfterExit=no

[Install]
WantedBy=default.target

