# machines/wsa/fsws-arch-config/systemd/user/thermalright-display.service
[Unit]
Description=Thermalright LCD Display Controller (NEXUS)
After=default.target

[Service]
Type=simple
ExecStart=%h/fsws-arch-config/display-control.sh start
ExecStop=%h/fsws-arch-config/display-control.sh stop
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target

