# machines/wsa/fsws-arch-config/systemd/weekly-system-update.timer
[Unit]
Description=Run weekly system update

[Timer]
# Runs once a week; Persistent ensures a missed run is caught up at next boot.
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target

