# configs/systemd/user/backup.service
[Unit]
Description=Rsync backup over Tailscale
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/backup.sh
