# configs/systemd/user/cache-cleanup.service
[Unit]
Description=Application cache cleanup

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/cache-cleanup.sh
